# docker/crontab
# Environment setup for cron jobs
PATH=/home/agent/.local/bin:/usr/local/bin:/usr/bin:/bin:/app/venv/bin
SHELL=/bin/bash

# Singularity Agent Cron Schedule
# All times are in container timezone (default: UTC)

# Hourly heartbeat - enqueue a cron run via the queue API
0 * * * * curl -s -X POST http://localhost:3001/api/queue/enqueue -H 'Content-Type: application/json' -d '{"type":"cron"}' >> /app/logs/heartbeat.log 2>&1

# Hourly vector index rebuild at :50
50 * * * * curl -s -X POST http://vector:5000/index >> /app/logs/vector-index.log 2>&1
