# docker/crontab
# Environment setup for cron jobs
PATH=/home/agent/.local/bin:/usr/local/bin:/usr/bin:/bin:/app/venv/bin
SHELL=/bin/bash

# Singularity Agent Cron Schedule
# All times are in container timezone (default: UTC)

# Hourly heartbeat - runs at the top of each hour
0 * * * * /app/scripts/run-agent.sh >> /app/logs/heartbeat.log 2>&1
